# src/examples/service/cron-schedule/cron.hourly/0100-full-backup.sh
#!/bin/bash

HOUR=$(date +%H)

if [[ "$HOUR" = "01" ]]; then
	cd /app
	./docker-backup-cron.sh
fi

